# upstart/smsd.conf
pre-start script
    mkdir -p /var/log/eri
end script
respawn
respawn limit 15 5
start on runlevel [2345]
stop on runlevel [06]
script
    cd /opt/smsd/
    python2 /opt/smsd/server.py 2>&1 >> /var/log/eri/smsd.log
end script
